# supervisor/kafka.conf
[program:kafka]
command=/usr/local/bin/kafka-start.sh
autostart=true
autorestart=true
